# meduim-lab/files/fauxpositif.service
[Unit]
Description=Faux Positif SSH Script

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/fauxpositif.sh
User=root
