# meduim-lab/files/legit_ssh.service
[Unit]
Description=Legitimate SSH Connection Service

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/legit_ssh.sh
User=root
